# telegram_music_bot/musicbot.service
[Unit]
Description=Telegram Music Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/telegram_music_bot
EnvironmentFile=/opt/telegram_music_bot/.env
ExecStart=/opt/telegram_music_bot/.venv/bin/python -m bot.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
